# user/protonmail.service
[Unit]
Description=Protonmail Bridge

[Service]
ExecStart=/usr/bin/proton-mail-bridge
Type=exec

[Install]
WantedBy=default.target
