# user/pulseaudio.service
[Unit]
Description=Pulseaudio server
Requires=jack-audio.service
After=jack-audio.service

[Service]
ExecStart=/usr/bin/pulseaudio
Type=exec

